# files/system/usr/lib/systemd/system/adw-gtk3-firstboot.service
[Unit]
Description=Flatpak GTK overrides and masks on first boot
After=flatpak-firstboot.service
Requires=flatpak-firstboot.service
ConditionPathExists=!/var/lib/silvercachy/.adw-gtk3-firstboot.done

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak override --system --filesystem=xdg-data/themes
ExecStart=/usr/bin/flatpak mask org.gtk.Gtk3theme.adw-gtk3
ExecStart=/usr/bin/flatpak mask org.gtk.Gtk3theme.adw-gtk3-dark
ExecStartPost=/usr/bin/mkdir -p /var/lib/silvercachy
ExecStartPost=/usr/bin/touch /var/lib/silvercachy/.adw-gtk3-firstboot.done
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target